# ansible/yum-reposync/swng-reposync.service.j2
[Unit]
Description=Sync SWNG Repositories with reposync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
{% for repo in swng_repos %}
{% if repo.enabled %}
{% if repo.module_platform_id is defined %}
ExecStart=/usr/bin/reposync -p {{ swng_mirror_path }}/ --repo {{ repo.name }} --setopt=module_platform_id={{ repo.module_platform_id }}
{% else %}
ExecStart=/usr/bin/reposync -p {{ swng_mirror_path }}/ --repo {{ repo.name }}
{% endif %}
ExecStartPost=/usr/bin/createrepo {{ swng_mirror_path }}/{{ repo.name }}/
{% endif %}
{% endfor %}
ExecStartPost=/usr/bin/python3 /opt/healthcheck/healthcheck_update.py --service sync --field repo --value {{ healthcheck_source_name }} --status OK
StandardOutput=append:{{ sync_log_file }}
StandardError=append:{{ sync_log_file }}
