# ansible/specific-version-rsync(Recomended)/swng-version-mirror.timer.j2
[Unit]
Description=Run CloudLinux {{ cloudlinux_version }} SWNG Mirror Sync Every {{ sync_interval_hours }} Hours
Requires={{ service_name }}.service

[Timer]
OnCalendar={{ timer_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
